# scripts/com.pilotpoint.county-discovery-health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pilotpoint.county-discovery-health</string>
  <key>ProgramArguments</key>
  <array>
    <string>PROJECT_DIR/.venv/bin/python</string>
    <string>PROJECT_DIR/county_discovery_health.py</string>
    <string>--output</string>
    <string>PROJECT_DIR/data/health/county-discovery-latest.json</string>
    <string>--strict-discrepancies</string>
    <string>--all-samples</string>
    <string>--notify-pcs</string>
    <string>--cleanup-orphans</string>
  </array>
  <key>WorkingDirectory</key>
  <string>PROJECT_DIR</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>6</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>PROJECT_DIR/data/health/county-discovery.log</string>
  <key>StandardErrorPath</key>
  <string>PROJECT_DIR/data/health/county-discovery-error.log</string>
</dict>
</plist>
